# files/spacewalk.conf
[supervisord]
nodaemon=true

[program:spacewalk]
command=/opt/spacewalk.sh
autorestart=false
